# install_crontab.sh.crontab
0,30 * * * * /bin/bash -l -c 'cd /home/deploy/house_to_trello && ruby main.rb'